# deploy/trump-trader-telegram.service
[Unit]
Description=Trump Trader Telegram Handler
After=network.target postgresql.service trump-trader.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trump_trader
Environment="PATH=/home/ubuntu/trump_trader/venv/bin:/usr/bin:/bin"
Environment="APP_ENV=production"
ExecStart=/home/ubuntu/trump_trader/venv/bin/python /home/ubuntu/trump_trader/telegram_bot_handler.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/trump_trader/logs/telegram.log
StandardError=append:/home/ubuntu/trump_trader/logs/telegram_error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

